# 0099FF-stack/worker/flower.conf
[program:flower]
command=/usr/local/opt/python/bin/flower --basic_auth=%(ENV_FLOWER_USER_NAME)s:%(ENV_FLOWER_PASSWORD)s
directory=/data/deploy/current
autorestart=true
