# 0099FF-stack/worker/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /etc/redis/redis.conf
autorestart=true 
priority=111
