# 0099FF-stack/worker/supervisor-djangorq-worker.conf
[program:djangorq-worker]
command=python manage.py rqworker default
autorestart=false
priority=999
stdout_logfile=syslog
stderr_logfile=syslog
