# web/python-nginx/conf/supervisor-uwsgi.conf
[program:uwsgi]
command=/etc/supervisor/conf.d/run_uwsgi.sh
user=uwsgi
autorestart=true
priority=500
stdout_syslog=true
stderr_syslog=true
