# 66CCFF-stack/worker/conf/supervisor-celery.conf
; ==================================
;  celery worker supervisor example
; ==================================

[program:celery]
; Set full path to celery program if using virtualenv
command=/etc/supervisor/conf.d/run_celery.sh 
directory=/data/deploy/current
numprocs=1
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; Supervisor struggles with killing celery workers. Next two settings appear to best combo
; (from https://github.com/Supervisor/supervisor/issues/721#issuecomment-184862648)
stopasgroup=true
stopsignal=INT

priority=999

stdout_logfile=syslog
stderr_logfile=syslog
