# 1213d64-stack/web/conf/supervisor-statsd.conf
[program:statsd]
command = /usr/bin/node /opt/statsd/stats.js  /opt/statsd/config_udp.js
stdout_logfile = syslog
stderr_logfile = syslog
autorestart = true