# 1213d64-stack/web/conf/supervisor-grafana.conf
[program:grafana-webapp]
directory = /opt/grafana/
command = /etc/supervisor/conf.d/run_grafana.sh
stdout_logfile = syslog
stderr_logfile = syslog
autorestart = true