# 1213d64-stack/web/conf/supervisor-carbon.conf
[program:carbon-cache]
command = /opt/graphite/bin/carbon-cache.py --pidfile /var/run/carbon-cache-a.pid --debug start
stdout_logfile = syslog
stderr_logfile = syslog
autorestart = true