# db/redis/conf/supervisor-redis.conf
[program:redis]
command=/usr/local/bin/redis-server /etc/redis/redis.conf
autorestart=true 
priority=111
user=redis

stdout_syslog=true
stderr_syslog=true