# db/postgres/conf/postgres.conf
[program:postgres]
command=/usr/lib/postgresql/9.4/bin/postgres -D /media/data/postgres/db -c config_file=%(ENV_POSTGRES_CONFIG_FILE)s
user=postgres
numprocs=1
autostart=true
autorestart=true

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

priority=5
