# db/postgres/conf/rsyslogd.conf
[program:rsyslog]
command = rsyslogd -n
autostart=true
autorestart=true
priority=1
