# 1213d65-stack/web/conf/supervisor-nginx.conf
[program:nginx]
command=/usr/local/nginx/sbin/nginx -c /usr/local/nginx/conf/nginx.conf
autorestart=true
priority=999
stdout_logfile=None
stderr_logfile=None
stdout_syslog=True
stderr_syslog=True

