# 1213d65-stack/web/conf/supervisor-export-datasources-and-dashboards.conf
[program:export-datasources-and-dashboards]
directory = /src
command = /bin/bash export-datasources-and-dashboards.sh
stdout_logfile = syslog
stderr_logfile = syslog
exitcodes = 0
autorestart = unexpected