# 0099FF-stack/web/uwsgi.conf
[program:uwsgi]
command=/usr/local/opt/python/bin/uwsgi /data/deploy/current/uwsgi.ini 
autorestart=true
priority=500
