# web/python-apache/nginx.conf
[program:nginx]
command=/usr/local/nginx/sbin/nginx
autorestart=true
priority=999
