# web/python-apache/uwsgi.conf
[program:uwsgi]
command=env -i /usr/local/opt/python/bin/uwsgi /data/deploy/current/uwsgi.ini
user=uwsgi
autorestart=true
priority=500
